# people/greg/splack-11.0/source/kde/kdebase/profile.d/kde.sh
#!/bin/sh
# KDE additions:
KDEDIR=/opt/kde
PATH=$PATH:$KDEDIR/bin
MANPATH="${MANPATH}:/opt/kde/man"
export KDEDIR PATH MANPATH
