# splack-8.0/source/kde/kdebase/profile.d/kde.csh
#!/bin/csh
# KDE additions:
if ( ! $?KDEDIR ) then
    setenv KDEDIR /opt/kde
endif
if ( $?PATH ) then
    setenv PATH $KDEDIR/bin:$PATH
else
    setenv PATH $KDEDIR/bin
endif
