# splack-8.0/source/kde/kdebase/profile.d/kde.sh
#!/bin/sh
# KDE additions:
KDEDIR=/opt/kde
PATH=$PATH:$KDEDIR/bin
export KDEDIR PATH
