# splack-8.0/source/kde/htdig/profile.d/htdig.csh
#!/bin/csh
# Environment path variables for the htdig package:
setenv PATH ${PATH}:/opt/www/htdig/bin
