# people/greg/splack-11.0/source/n/htdig/profile.d/htdig.sh
#!/bin/sh
# Environment variables for the htdig package:
PATH="$PATH:/opt/www/htdig/bin"
export PATH
